# data/libertine-lxc-manager.conf
description "Service to manage Libertine LXC's"

start on started unity8
stop on desktop-end

respawn

exec /usr/bin/libertine-lxc-manager
