# data/libertine.service
[Unit]
Description=Libertine Service

[Service]
Type=dbus
BusName=com.canonical.libertine.Service
ExecStart=/usr/bin/libertined
