# data/libertine-lxc-manager.conf
description "Service to manage Libertine LXC's"
author "Christopher Townsend <christopher.townsend@canonical.com>"

start on started unity8
stop on desktop-end

respawn

pre-start script
	dpkg -s python3-libertine-lxc > /dev/null 2>&1 || { stop; exit 0; } 
end script

exec /usr/bin/libertine-lxc-manager
