# data/puritine-click.conf
description "Puritine Click chroot linking"

start on starting unity8

script
PACKAGE_PATH=`click pkgdir com.ubuntu.puritine`
CONTAINER_NAME=puritine
CHROOT_PATH=$PACKAGE_PATH/libertine-data/libertine-container/$CONTAINER_NAME

if [ -x $CHROOT_PATH/rootfs ] ; then
	# Link the chroot
	if [ ! -L $HOME/.cache/libertine-container/$CONTAINER_NAME/rootfs ] ; then
		mkdir -p $HOME/.cache/libertine-container/$CONTAINER_NAME/
		ln -s $CHROOT_PATH/rootfs $HOME/.cache/libertine-container/$CONTAINER_NAME/rootfs
	fi

	# Copy or merge the container config files
	if [ ! -e $HOME/.local/share/libertine/ContainersConfig.json ] ; then
		mkdir -p $HOME/.local/share/libertine/
		cp $PACKAGE_PATH/libertine-config/libertine/ContainersConfig.json $HOME/.local/share/libertine/ContainersConfig.json
	elif [ -L $HOME/.local/share/libertine/ContainersConfig.json ] ; then
		rm $HOME/.local/share/libertine/ContainersConfig.json
		cp $PACKAGE_PATH/libertine-config/libertine/ContainersConfig.json $HOME/.local/share/libertine/ContainersConfig.json
	else
		libertine-container-manager merge-configs -f $PACKAGE_PATH/libertine-config/libertine/ContainersConfig.json    
	fi

	# Create and copy the user-data dir from the click package
	if [ ! -d $HOME/.local/share/libertine-container/user-data/$CONTAINER_NAME ] ; then
		cp -dR $PACKAGE_PATH/libertine-config/libertine-container/ $HOME/.local/share/libertine-container/
	fi
else
	rm -rf $HOME/.cache/libertine-container/$CONTAINER_NAME
	if [ $(libertine-container-manager list | grep -v ${CONTAINER_NAME} | wc -l) -eq 0 ]; then
		rm -rf $HOME/.local/share/libertine
	fi
	rm -rf $HOME/.local/share/libertine-container/user-data/$CONTAINER_NAME/.config/
fi
end script
