# data/puritine-click.conf
description "Puritine Click chroot linking"

start on started unity8

script
PACKAGE_PATH=`click pkgdir com.ubuntu.puritine`
CONTAINER_NAME=puritine
CHROOT_PATH=$PACKAGE_PATH/libertine-data/libertine-container/$CONTAINER_NAME

if [ -x $CHROOT_PATH/rootfs ] ; then
	# Link the chroot
	if [ ! -L $HOME/.cache/libertine-container/$CONTAINER_NAME/rootfs ] ; then
		mkdir -p $HOME/.cache/libertine-container/$CONTAINER_NAME/
		ln -s $CHROOT_PATH/rootfs $HOME/.cache/libertine-container/$CONTAINER_NAME/rootfs
	fi

	# Link the container config
	if [ ! -L $HOME/.local/share/libertine/ContainersConfig.json ] ; then
		mkdir -p $HOME/.local/share/libertine/
		ln -s $PACKAGE_PATH/libertine-config/libertine/ContainersConfig.json $HOME/.local/share/libertine/ContainersConfig.json
	fi
else
	rm -rf $HOME/.cache/libertine-container/$CONTAINER_NAME
	rm -rf $HOME/.local/share/libertine
	rm -rf $HOME/.local/share/libertine-container/user-data/$CONTAINER_NAME/.config/
fi
end script
