# data/replace-libertine-scope.conf
description "Remove Libertine scope click"
author "Dalton Durst <dalton@ubports.com>"

# This job removes the Libertine click. This effectively replaces the scope
# since it is now included in the image.

start on starting unity8

script
    echo ""
    if CLICK=`click info libertine-scope.ubuntu`; then
        GET_VERSION="import json,sys; json_doc=json.loads(sys.stdin.read()); print(json_doc['version'])"
        echo "Replacing Libertine scope"
        LIBERTINE_SCOPE_VERSION=`echo $CLICK | python3 -c "$GET_VERSION"`
        pkcon remove "libertine-scope.ubuntu;$LIBERTINE_SCOPE_VERSION;all;local:click"
    else
        echo "Scope click not installed, no need to replace."
    fi
end script
